# rootfs/etc/cron.d/notebooks
0 0 * * * /user/bin/bash /root/update_notebooks.sh
